# files/supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/sshd.out.log
stderr_logfile=/var/log/supervisor/sshd.err.log